# webrtc-body.service
[Unit]
Description=Body WebRTC Server
After=network.target

[Service]
Type=simple
User=body
WorkingDirectory=/home/body/webrtc-body
ExecStart=/home/body/webrtc-body/server.py
Restart=always
Environment=DISPLAY=:0

[Install]
WantedBy=graphical.target
